# rns-icn.service
[Unit]
Description=ICN Origin Server (rns-icn)
# Rides the shared rnsd transport daemon; start after it is up.
After=network-online.target rnsd.service
Wants=network-online.target
Requires=rnsd.service

[Service]
Type=simple
User=root
Group=root
Environment="HOME=/root"
ExecStart=/opt/rns-icn-venv/bin/icn-server --config /etc/rns-icn/icn.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
