# rnsd.service
[Unit]
Description=Reticulum Transport Daemon — rns-icn ICN network
Documentation=https://hermes-agent.nousresearch.com/docs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/rnsd --config /etc/rnsd
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
